# crontabs/www-data
*/10 * * * *  cd /var/www/html && wp cron event run --due-now
