# crontabs/www-data
*/10 * * * *  wp cron event run --due-now
